# api/supervisor.ini
[program:boswps]
directory = /root/referendum/api
command = gunicorn -w 4 api:app
stderr_logfile = /root/referendum/api/supervisord.d/log/boswps-stderr.log
stdout_logfile = /root/referendum/api/supervisord.d/log/boswps-stdout.log
autostart = false
autorestart = true
user = root
stopasgroup=true
killasgroup=true