# api/supervisor.ini
[program:boswps]
directory = /path/to/boswps-api
command = gunicorn -w 4 api:app
stderr_logfile = /path/to/supervisord.d/log/boswps-stderr.log
stdout_logfile = /path/to/supervisord.d/log/boswps-stdout.log
autostart = false
autorestart = true
user = root
stopasgroup=true
killasgroup=true
